# deploy/gastro-frontend.service
[Unit]
Description=Gastro Frontend (http.server Port 8080)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bestellung
Group=bestellung
WorkingDirectory=/opt/gastro-system
ExecStart=/opt/gastro-system/venv/bin/python /opt/gastro-system/deploy/frontend_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
