# deploy/gastro-wifi-autoswitch.service
[Unit]
Description=Gastro WiFi auto-switch (best signal)
After=network.target NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/gastro-system/deploy/wifi-autoswitch.sh

